# manifests/modules/workhammer/templates/workhammer.conf.erb
[program:workhammer]
command=/opt/workhammer/bin/workhammer
numprocs=1
user=workhammer
stdout_logfile=/var/log/supervisor/workhammer.log
stderr_logfile=/var/log/supervisor/workhammer.err.log
autostart=true
autorestart=true
